# autostart/gpg-connect-agent.desktop
[Desktop Entry]
Exec=gpg-connect-agent reloadagent /bye
Icon=
Name=gpg-connect-agent
Path=
Terminal=False
Type=Application
